# Task_3/ansible-staging-app/roles/webserver/templates/my_node_app.service.j2
[Unit]
Description=My Node.js Application
After=network.target

[Service]
User={{ ansible_user }}
WorkingDirectory={{ app_dir }}
ExecStart=/usr/bin/node {{ app_dir }}/app.js
Restart=always

[Install]
WantedBy=multi-user.target